# com.citapadron.checker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.citapadron.checker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/satorrea/PycharmProjects/cita-padron-madrid/.venv/bin/python</string>
        <string>/Users/satorrea/PycharmProjects/cita-padron-madrid/run.py</string>
        <string>--once</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/satorrea/PycharmProjects/cita-padron-madrid</string>

    <key>StartInterval</key>
    <integer>1800</integer>  <!-- Cada 30 minutos (1800 segundos) -->

    <key>StandardOutPath</key>
    <string>/Users/satorrea/PycharmProjects/cita-padron-madrid/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/satorrea/PycharmProjects/cita-padron-madrid/launchd_error.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
